# dot_config/autostart/wallpaper-monitor.desktop
[Desktop Entry]
Type=Application
Name=Wallpaper Monitor
Comment=ディスプレイ接続時に壁紙を自動設定
Exec=$HOME/.local/bin/wallpaper-monitor.sh
Hidden=false
NoDisplay=true
X-GNOME-Autostart-enabled=true
